# config/common/ripe-atlas.service.in
[Unit]
Description=Atlas Probe
After=network-online.target syslog.target

[Service]
User=@ripe_atlas_user@
Group=@ripe_atlas_group@
Environment=HOME=@atlas_spooldir@
WorkingDirectory=@atlas_spooldir@
ExecStart=@sbindir@/ripe-atlas
ExecStop=/usr/bin/killall -n 0 telnetd perd eperd eooqd
Restart=always
TimeoutStopSec=60
DevicePolicy=closed
KeyringMode=private
ProtectHome=yes
PrivateTmp=yes
RemoveIPC=yes
CapabilityBoundingSet=~CAP_WAKE_ALARM
CapabilityBoundingSet=~CAP_SYS_BOOT
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
CapabilityBoundingSet=~CAP_IPC_LOCK
CapabilityBoundingSet=~CAP_SYS_CHROOT
CapabilityBoundingSet=~CAP_SYS_TIME
CapabilityBoundingSet=~CAP_KILL
CapabilityBoundingSet=~CAP_BLOCK_SUSPEND
CapabilityBoundingSet=~CAP_SYS_PTRACE
CapabilityBoundingSet=~CAP_AUDIT_CONTROL
CapabilityBoundingSet=~CAP_MKNOD

[Install]
WantedBy=multi-user.target
